# celery.conf
[program:celery]
directory = /code/sandbox
command = celery -A sandbox worker -l info
stdout_logfile = /var/log/celeryd.log
stderr_logfile = /var/log/celeryd.err